# etc/cron.d/firewall-on-reboot
@reboot root /etc/lighttpd/returnco-de-firewall.sh

